# deploy/experiment-daemon.service
[Unit]
Description=Quantum Inspire Experiment Daemon
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/quantuminspire
ExecStart=/opt/quantuminspire/.venv/bin/python agents/experiment_daemon.py --daemon --interval 300
Restart=on-failure
RestartSec=60
StandardOutput=journal
StandardError=journal

Environment=HOME=/root
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
